# feh-randomize.service
[Unit]
Description=wallpaper

[Service]
ExecStart=/usr/bin/fen-randomize

[Install]
WantedBy=multi-user.target
