# feh-randomize.timer
[Unit]
Description=Executa foo semanalmente e na inicialização

[Timer]
OnBootSec=1in
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
